# docker-flask/docker-app/nginflask.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=/usr/local/bin/uwsgi --ini /flaskApp/uwsgi.ini

[program:nginx]
command=/usr/sbin/nginx